# itamae/hosts/sumi/files/etc/systemd/system/conserver.service
[Unit]
Description = console server

[Service]
ExecStart = /usr/bin/conserver
Restart = always
Type = simple

[Install]
WantedBy = multi-user.target
